# systemd/browser-agent-observer.target
[Unit]
Description=browser-agent-observer (backend + proxy + screencast)

[Install]
WantedBy=default.target
